# systemd/riemann-dash.service
[Unit]
Description=Riemann Dash
Requires=riemann.service
After=riemann.service

[Service]
User=ubuntu
ExecStart=/usr/bin/riemann-dash
